# docker_prism/dev/etc/supervisord/supervisord.d/apache.ini
[program:apache]
command=/usr/sbin/httpd -D FOREGROUND      ; Pour RHEL/Rocky/CentOS
autostart=true
autorestart=true
stdout_logfile=/var/log/apache.log
stderr_logfile=/var/log/apache-error.log
user=root
